# lcd_display.service
[Unit]
Description=LCD Library

[Service]
ExecStart=/usr/local/bin/display
Restart=always

[Install]
WantedBy=multi-user.target